# services/forgejo/renovate.timer
[Unit]
Description=Run Renovate bot

[Timer]
OnCalendar=hourly
Persistent=true
Unit=renovate.service

[Install]
WantedBy=timers.target