# services/inventree/inventree-redis.container
[Unit]
Description=InvenTree cache

[Container]
Image=docker.io/valkey/valkey:8.1-alpine
AutoUpdate=registry
Pod=inventree.pod
# Healthcheck
HealthCmd="redis-cli ping || exit 1"
Notify=healthy

[Service]
Restart=on-failure
