# services/traefik/imaps.socket
[Socket]
ListenStream=993
FileDescriptorName=imaps
Service=traefik.service

[Install]
WantedBy=sockets.target
