# services/traefik/ldap.socket
[Socket]
ListenStream=636
FileDescriptorName=ldapsecure
Service=traefik.service

[Install]
WantedBy=sockets.target
