# services/traefik/smtp.socket
[Socket]
ListenStream=25
FileDescriptorName=smtp
Service=traefik.service

[Install]
WantedBy=sockets.target
