# services/traefik/smtps.socket
[Socket]
ListenStream=465
FileDescriptorName=smtps
Service=traefik.service

[Install]
WantedBy=sockets.target
