# services/traefik/ssh.socket
[Socket]
ListenStream=2222
FileDescriptorName=sshalt
Service=traefik.service

[Install]
WantedBy=sockets.target
